# supervisor.conf
[program:yp]
command=/var/www/yprogrammer/gunicorn_start
user=root
autostart=true
autorestart=true